# systemd/goes-sync.timer
[Unit]
Description=Run the GOES product sync hourly

[Timer]
# Hourly keeps the SD-card backlog small (~1.3 GB) instead of a full day's ~30 GB,
# so a single failed run can't fill the card and stall the decode. The 5-min
# settle window in goes-sync.sh means each run still skips in-flight products.
OnCalendar=*-*-* *:00:00
Persistent=true
# Small jitter so the run doesn't land exactly on the hour with other cron-like load.
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
